# deploy/ML-Netease_Url.service
[Unit]
Description=ML Netease URL Gunicorn service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/ML-Netease_url
EnvironmentFile=/etc/ML-Netease_url.env
ExecStart=/opt/ML-Netease_url/.venv/bin/gunicorn -c /opt/ML-Netease_url/deploy/gunicorn.conf.py wsgi:app
Restart=always
RestartSec=5
KillSignal=SIGQUIT
TimeoutStopSec=30
PrivateTmp=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
